# bash/.config/bash/profile.d/alias.sh
#!/bin/bash

alias reload='. ~/.bashrc;'
