# bw/.config/systemd/user/bw-expire.service
[Unit]
Description=Run bw-expire

[Service]
ExecStart=bash %h/.local/bin/bw-expire.sh
